# form-handler/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:form-handler]
command=/opt/swimwithjj/form-handler
